# skald-server.service
[Unit]
Description=Skald-Go Transcriber Service
After=network.target

[Service]
Type=simple
User=%u
WorkingDirectory=%h/Documents/Dev/skald-go
ExecStart=%h/Documents/Dev/skald-go/run-server.sh
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=default.target 